# conf/aeneria.cron
*/20 * * * * #USER# php7.3 #DESTDIR#/bin/console aeneria:fetch-data
